# docker/python/supervisord.conf
[program:miner]
command=/bin/bash -c "/usr/local/bin/python /app/start_miner.py"
numprocs=3
process_name=%(program_name)-%(process_num)s
stdout_logfile=/app/logs/miner.log
redirect_stdout=true